# loop/com.imaginarycreatures.loop.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.imaginarycreatures.loop</string>

  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>REPLACE_WITH_ABSOLUTE_PATH/loop/run-loop.sh</string>
  </array>

  <!-- Every 10 hours (36000 seconds). -->
  <key>StartInterval</key>
  <integer>36000</integer>

  <key>RunAtLoad</key>
  <false/>

  <key>StandardOutPath</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH/loop/launchd.out.log</string>
  <key>StandardErrorPath</key>
  <string>REPLACE_WITH_ABSOLUTE_PATH/loop/launchd.err.log</string>
</dict>
</plist>
